# systemd/foomuuri-resolve.service
[Unit]
Description=Multizone bidirectional nftables firewall - Resolve hostnames

[Service]
ExecStart=/usr/sbin/foomuuri resolve
