# systemd/foomuuri-iplist.service
[Unit]
Description=Multizone bidirectional nftables firewall - Update iplist entries

[Service]
ExecStart=/usr/sbin/foomuuri iplist
