# systemd/foomuuri-iplist.timer
[Unit]
Description=Multizone bidirectional nftables firewall - Refresh iplist entries
Documentation=https://foomuuri.foobar.fi/latest/
Requires=foomuuri.service
PartOf=foomuuri.service

[Timer]
OnActiveSec=3m
OnUnitInactiveSec=15m
