# systemd/foomuuri-resolve.timer
[Unit]
Description=Multizone bidirectional nftables firewall - Resolve hostnames
Documentation=https://github.com/FoobarOy/foomuuri/wiki
Requires=foomuuri.service
PartOf=foomuuri.service

[Timer]
OnBootSec=3m
OnUnitInactiveSec=15m
