# systemd/foomuuri-monitor.service
[Unit]
Description=Multizone bidirectional nftables firewall - Target monitor
Documentation=https://github.com/FoobarOy/foomuuri/wiki
After=network-online.target
After=foomuuri.service
PartOf=foomuuri.service

[Service]
Type=notify
ExecStart=/usr/sbin/foomuuri monitor
ExecReload=kill -HUP $MAINPID
